# services/proxy/deploy/systemd/veilkey-proxy-daluser.service
[Unit]
Description=VeilKey Proxy observer enforcement for daluser
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
ExecStart=/root/veilkey-proxy/bin/veilkey-proxy observe --uid 1008 --format json --only-suspicious --enforce-kill
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
