# services/proxy/deploy/host/veilkey-root-ai-egress-guard@.service
[Unit]
Description=VeilKey root-ai cgroup egress guard for profile %i
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Environment=VEILKEY_ROOT_AI_PROFILE=%i
ExecStart=/usr/local/bin/veilkey-root-ai-egress-guard %i apply
ExecStop=/usr/local/bin/veilkey-root-ai-egress-guard %i remove
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
