# services/proxy/deploy/lxc/veilkey-egress-proxy@.service
[Unit]
Description=VeilKey egress proxy (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
EnvironmentFile=-/etc/veilkey/proxy.env
StandardOutput=append:/var/log/veilkey-proxy/%i.jsonl
StandardError=append:/var/log/veilkey-proxy/%i.log
ExecStart=/usr/local/bin/veilkey-proxy-launch %i
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
